# debian/statsd.upstart
description "statsd"
author      "rdio"

start on startup
stop on shutdown

script
    # We found $HOME is needed. Without it, we ran into problems
    export HOME="/root"

    exec sudo -u fabric /usr/share/statsd/scripts/start
end script
