# debian/statsd.upstart
description "statsd"
author      "etsy"

start on startup
stop on shutdown

script
    chdir /usr/share/statsd

    exec sudo -u _statsd /usr/share/statsd/scripts/start
end script
